# deb/kaltura-sphinx/debian/prerm
#!/bin/bs
if [ "$1" = "remove" -o $1 = "purge" ] ; then
    /sbin/service kaltura-sphinx stop >/dev/null 2>&1
    update-rc.d -f kaltura-sphinx remove
fi
