# deb/kaltura-playkit-bundler/debian/preinst
#!/bin/bash
# we really only want this when upgrading but since we have to use '|| true' anyhow, it doesn't matter
invoke-rc.d kaltura-playkit-bundler stop || true
