# deb/kaltura-playkit-bundler/debian/prerm
#!/bin/sh
if [ "$1" = "remove" -o $1 = "purge" ] ; then
    /usr/sbin/service kaltura-playkit-bundler stop >/dev/null 2>&1
    update-rc.d -f kaltura-playkit-bundler remove
fi
